# americancoind.service
[Unit]
Description=Americancoin Daemon

[Service]
Type=forking
ExecStart=%h/bin/americancoind

[Install]
WantedBy=default.target
