# bitcoind.service
[Unit]
Description=Bitccoin Daemon

[Service]
Type=forking
ExecStart=/usr/local/bitcoin/bin/bitcoind

[Install]
WantedBy=default.target
