# clamd.service
[Unit]
Description=Clamcoin Daemon

[Service]
Type=forking
ExecStart=%h/bin/%N

[Install]
WantedBy=default.target
