# ipfs.service
[Unit]
Description=IPFS, InterPlanetary File System

[Service]
Type=simple
ExecStart=%h/bin/%N daemon
ExecStop=%h/bin/%N shutdown
Restart=always
RestartSec=60

[Install]
WantedBy=default.target
